# src/pi/PiDemos/PiDemo.BeaconsDemo/pibeacon.service
[Unit]
Description=Bluetooth iBeacon transmitter
After=bluetooth.service

[Service]
Type=simple
ExecStart=/opt/pibeacon/PiBleBeacon
Restart=always
RestartSec=2
# Optional: run as a non-root service account that owns /opt/pibeacon
User=pi
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
